# ansible/roles/serf/templates/supervisord-serf.conf.j2
[program:serf]
command=serf agent -tag role={{serf_role}} -event-handler=/openvswitch-handler.sh
numprocs=1
autostart=true
autorestart=true
